# config/supervisor.conf
[program:wox]
command = uwsgi /var/www/wox/config/uwsgi.ini
autorestart = true
stdout_logfile = /var/www/wox/log/uwsgi/out.log
redirect_stderr = true
