# config/supervisor.conf
[program:wox]
command = uwsgi /app/wox/config/uwsgi.ini
autorestart = true
stdout_logfile = /app/wox/log/uwsgi/out.log
redirect_stderr = true
